# package/files/etc/init.d/dns-smart-routing
#!/bin/sh /etc/rc.common

START=95
USE_PROCD=1

start_service() {
    # Ensure directories exist
    mkdir -p /etc/dns-smart-routing

    # Initialize default state file if missing or corrupt
    if [ ! -f /etc/dns-smart-routing/state.json ] || [ ! -s /etc/dns-smart-routing/state.json ] || ! jq -e . /etc/dns-smart-routing/state.json >/dev/null 2>&1; then
        echo '{"state":"CLEAN","consecutive_fail":0,"consecutive_ok":0,"last_change":0}' > /etc/dns-smart-routing/state.json
    fi

    # Ensure dynamic servers file exists for dnsmasq
    touch /tmp/dnsmasq_dynamic_servers.conf

    # Verify and add dnsmasq config line safely (no duplicates)
    if [ -f /etc/dnsmasq.conf ]; then
        sed -i '/servers-file=\/tmp\/dnsmasq_dynamic_servers.conf/d' /etc/dnsmasq.conf
        echo "servers-file=/tmp/dnsmasq_dynamic_servers.conf" >> /etc/dnsmasq.conf
        /etc/init.d/dnsmasq restart
    fi

    # Setup cron entry safely without duplicates
    if [ -f /etc/crontabs/root ]; then
        sed -i '/dns_smart_probe.sh/d' /etc/crontabs/root
        echo "*/1 * * * * /usr/bin/dns_smart_probe.sh && /usr/bin/dns_smart_apply.sh" >> /etc/crontabs/root
        /etc/init.d/cron restart
    fi
}

stop_service() {
    # Clean up cron entry safely
    if [ -f /etc/crontabs/root ]; then
        sed -i '/dns_smart_probe.sh/d' /etc/crontabs/root
        /etc/init.d/cron restart
    fi
}
