# package/files/etc/init.d/dns-smart-routing
#!/bin/sh /etc/rc.common

START=95
USE_PROCD=1

STATE_DIR="/etc/dns-smart-routing"

_dnsmasq_conf_add() {
    local line="servers-file=/tmp/dnsmasq_dynamic_servers.conf"
    local conf="/etc/dnsmasq.conf"
    [ ! -f "$conf" ] && return 1

    local count
    count=$(grep -c "$line" "$conf" 2>/dev/null || echo "0")
    if [ "$count" -eq 1 ]; then
        # Exactly one exists, no change needed
        return 0
    fi

    # Rewrite conf file: remove any lines containing servers-file=/tmp/dnsmasq_dynamic_servers.conf, then append exactly one
    local tmp="${conf}.tmp"
    grep -v "servers-file=/tmp/dnsmasq_dynamic_servers.conf" "$conf" > "$tmp" 2>/dev/null \
        && printf '%s\n' "$line" >> "$tmp" 2>/dev/null \
        && mv "$tmp" "$conf" 2>/dev/null \
        || { rm -f "$tmp" 2>/dev/null; return 1; }
    return 2 # 2 means modified
}

_cron_add() {
    local entry="*/1 * * * * /usr/bin/dns_smart_probe.sh && /usr/bin/dns_smart_apply.sh"
    local crontab="/etc/crontabs/root"
    touch "$crontab" 2>/dev/null

    local tmp="${crontab}.tmp"
    grep -v "dns_smart_probe.sh" "$crontab" > "$tmp" 2>/dev/null
    printf '%s\n' "$entry" >> "$tmp" 2>/dev/null

    if cmp -s "$tmp" "$crontab" 2>/dev/null; then
        rm -f "$tmp" 2>/dev/null
    else
        mv "$tmp" "$crontab" 2>/dev/null
        /etc/init.d/cron restart 2>/dev/null || true
    fi
}

_cron_remove() {
    local crontab="/etc/crontabs/root"
    [ ! -f "$crontab" ] && return

    local tmp="${crontab}.tmp"
    grep -v "dns_smart_probe.sh" "$crontab" > "$tmp" 2>/dev/null

    if cmp -s "$tmp" "$crontab" 2>/dev/null; then
        rm -f "$tmp" 2>/dev/null
    else
        mv "$tmp" "$crontab" 2>/dev/null
        /etc/init.d/cron restart 2>/dev/null || true
    fi
}

start_service() {
    mkdir -p "$STATE_DIR" 2>/dev/null

    # Clear temp counters
    rm -f /tmp/dns_smart_fail_count 2>/dev/null

    # Validate and repair state if missing, empty, invalid JSON, or not in NORMAL/FAILOVER
    local sf="$STATE_DIR/state.json"
    local is_corrupt=0
    if [ ! -f "$sf" ] || [ ! -s "$sf" ]; then
        is_corrupt=1
    else
        local state_val
        state_val=$(jq -r '.state' "$sf" 2>/dev/null)
        if [ $? -ne 0 ] || [ "$state_val" != "NORMAL" -a "$state_val" != "FAILOVER" ]; then
            is_corrupt=1
        fi
    fi

    if [ $is_corrupt -eq 1 ]; then
        printf '{"state":"NORMAL"}\n' > "$sf" 2>/dev/null
    fi

    touch /tmp/dnsmasq_dynamic_servers.conf 2>/dev/null

    _dnsmasq_conf_add
    local add_rc=$?
    if [ $add_rc -eq 2 ]; then
        /etc/init.d/dnsmasq restart 2>/dev/null || true
    fi

    _cron_add
}

stop_service() {
    _cron_remove
    rm -f /tmp/dns_smart_fail_count 2>/dev/null
}
